# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true

[program:reverb]
command=php artisan reverb:start
directory=/var/www
autostart=true
autorestart=true
